# systemd/helm-runner.service
[Unit]
Description=HELM Autonomous Model Runner Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/hoch_agent_swarm
ExecStart=/usr/bin/python3 -u /root/hoch_agent_swarm/scripts/helm_autonomy_runner.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
