# deploy/remote-relay/systemd/hoch-agent-swarm.service
[Unit]
Description=HOCH Agent Swarm / Prompt Brain Remote Control Plane
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/app
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down

[Install]
WantedBy=multi-user.target
