# hsf/deploy/com.hoch.model-upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.hoch.model-upgrade</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string><string>-lc</string>
    <string>cd ~/hoch_agent_swarm && python3 scripts/model_upgrade.py</string>
  </array>
  <key>StartCalendarInterval</key><dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/tmp/hoch-model-upgrade.log</string>
  <key>StandardErrorPath</key><string>/tmp/hoch-model-upgrade.log</string>
  <key>RunAtLoad</key><false/>
</dict></plist>
