# has_live_project_tracker/tracker-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hochster71.has.tracker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>server.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/michaelhoch/hoch_agent_swarm/has_live_project_tracker</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TRACKER_PORT</key>
        <string>3001</string>
        <key>UI_USER</key>
        <string>admin</string>
        <key>UI_PASS</key>
        <string>change-this-password</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/michaelhoch/hoch_agent_swarm/has_live_project_tracker/tracker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/michaelhoch/hoch_agent_swarm/has_live_project_tracker/tracker.log</string>
</dict>
</plist>
